# nvme-temp.service
[Unit]
Description=NVMe0 temperature monitor daemon

[Service]
ExecStart=/usr/local/bin/nvme0-temp

[Install]
WantedBy=multi-user.target
