# pkg-list.service
[Service]
ExecStart=/etc/systemd/system/pkg-list.sh

[Install]
WantedBy=default.target
