# services/sync-pkg-list.timer
[Unit]
Description=Sync packages hourly and on boot

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
