# services/timeshift-archive.timer
[Unit]
Description=Run timeshift-archive service every month

[Timer]
OnCalendar=*-*-01 00:00:00

[Install]
WantedBy=timers.target
