# services/timeshift-archive.timer
[Unit]
Description=Run timeshift-archive service every 7 days

[Timer]
OnCalendar=Mon *-*-* 00:00:00

[Install]
WantedBy=timers.target
